# deploy/sspanel-hy2-port-sync.timer
[Unit]
Description=Periodically sync SSPanel-UIM Hysteria2 node port

[Timer]
OnBootSec=30s
OnUnitActiveSec=30s
AccuracySec=5s
Persistent=true
Unit=sspanel-hy2-port-sync.service

[Install]
WantedBy=timers.target
