# deploy/sspanel-hy2-adapter.service
[Unit]
Description=SSPanel-UIM Hysteria 2 HTTP Auth Adapter
After=network-online.target hysteria-server.service
Wants=network-online.target

[Service]
Type=simple
User=sspanel-hy2
Group=sspanel-hy2
ExecStart=/usr/local/bin/sspanel-hy2-adapter -config /etc/sspanel-hy2-adapter/config.yaml
Restart=on-failure
RestartSec=3s
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
StateDirectory=sspanel-hy2-adapter
ReadWritePaths=/var/lib/sspanel-hy2-adapter

[Install]
WantedBy=multi-user.target
