# deploy/sspanel-hy2-port-sync.service
[Unit]
Description=Sync SSPanel-UIM Hysteria2 node port to Docker Compose
After=docker.service network-online.target
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/sspanel-hy2-adapter
ExecStart=/opt/sspanel-hy2-adapter/scripts/sync-panel-port.sh
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
NoNewPrivileges=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/opt/sspanel-hy2-adapter /run/docker.sock
